# Firmware/usr/lib/opkg/info/boot-play.list
/usr/res/boot-play/part0/pic_011.jpg
/usr/res/boot-play/part1/pic_132.jpg
/usr/res/boot-play/part1/pic_213.jpg
/usr/res/boot-play/part0/pic_054.jpg
/usr/res/boot-play/part0/pic_019.jpg
/usr/res/boot-play/part1/pic_094.jpg
/usr/res/boot-play/part1/pic_175.jpg
/usr/res/boot-play/part0/pic_020.jpg
/usr/res/boot-play/part1/pic_141.jpg
/usr/res/boot-play/part1/pic_106.jpg
/usr/res/boot-play/part0/pic_063.jpg
/usr/res/boot-play/part0/pic_028.jpg
/usr/res/boot-play/part1/pic_184.jpg
/usr/res/boot-play/part1/pic_149.jpg
/usr/res/boot-play/desc.txt
/usr/res/boot-play/part1/pic_150.jpg
/usr/res/boot-play/part1/pic_115.jpg
/usr/res/boot-play/part0/pic_072.jpg
/usr/res/boot-play/part0/pic_037.jpg
/usr/res/boot-play/part1/pic_193.jpg
/usr/res/boot-play/part1/pic_158.jpg
/usr/res/boot-play/part0/pic_003.jpg
/usr/res/boot-play/part1/pic_124.jpg
/usr/res/boot-play/part1/pic_205.jpg
/usr/res/boot-play/part0/pic_081.jpg
/usr/res/boot-play/part0/pic_046.jpg
/usr/res/boot-play/part1/pic_167.jpg
/usr/res/boot-play/part0/pic_012.jpg
/usr/res/boot-play/part1/pic_133.jpg
/usr/res/boot-play/part1/pic_214.jpg
/usr/res/boot-play/part0/pic_055.jpg
/usr/res/boot-play/part1/pic_095.jpg
/usr/res/boot-play/part1/pic_176.jpg
/usr/res/boot-play/part0/pic_021.jpg
/usr/res/boot-play/part1/pic_142.jpg
/usr/res/boot-play/part1/pic_107.jpg
/usr/res/boot-play/part0/pic_064.jpg
/usr/res/boot-play/part0/pic_029.jpg
/usr/res/boot-play/part1/pic_185.jpg
/usr/res/boot-play/part0/pic_030.jpg
/usr/res/boot-play/part1/pic_151.jpg
/usr/res/boot-play/part1/pic_116.jpg
/usr/res/boot-play/part0/pic_073.jpg
/usr/res/boot-play/part0/pic_038.jpg
/usr/res/boot-play/part1/pic_194.jpg
/usr/res/boot-play/part1/pic_159.jpg
/usr/res/boot-play/part0/pic_004.jpg
/usr/res/boot-play/part1/pic_160.jpg
/usr/res/boot-play/part1/pic_125.jpg
/usr/res/boot-play/part1/pic_206.jpg
/usr/res/boot-play/part0/pic_082.jpg
/usr/res/boot-play/part0/pic_047.jpg
/usr/res/boot-play/part1/pic_168.jpg
/usr/res/boot-play/boot.mp3
/usr/res/boot-play/part0/pic_013.jpg
/usr/res/boot-play/part1/pic_134.jpg
/usr/res/boot-play/part0/pic_056.jpg
/usr/res/boot-play/part1/pic_096.jpg
/usr/res/boot-play/part1/pic_177.jpg
/usr/res/boot-play/part1/pic_100.jpg
/usr/res/boot-play/part0/pic_022.jpg
/usr/res/boot-play/part1/pic_143.jpg
/usr/res/boot-play/part1/pic_108.jpg
/usr/res/boot-play/part0/pic_065.jpg
/usr/res/boot-play/part1/pic_186.jpg
/usr/res/boot-play/part0/pic_031.jpg
/usr/res/boot-play/part1/pic_152.jpg
/usr/res/boot-play/part1/pic_117.jpg
/usr/res/boot-play/part0/pic_074.jpg
/usr/res/boot-play/part0/pic_039.jpg
/usr/res/boot-play/part1/pic_195.jpg
/sbin/boot-play
/usr/res/boot-play/part0/pic_040.jpg
/usr/res/boot-play/part0/pic_005.jpg
/usr/res/boot-play/part1/pic_161.jpg
/usr/res/boot-play/part1/pic_126.jpg
/usr/res/boot-play/part1/pic_207.jpg
/usr/res/boot-play/part0/pic_083.jpg
/usr/res/boot-play/part0/pic_048.jpg
/usr/res/boot-play/part1/pic_169.jpg
/usr/res/boot-play/part0/pic_014.jpg
/usr/res/boot-play/part1/pic_170.jpg
/usr/res/boot-play/part1/pic_135.jpg
/usr/res/boot-play/part0/pic_057.jpg
/usr/res/boot-play/part1/pic_097.jpg
/usr/res/boot-play/part1/pic_178.jpg
/usr/res/boot-play/part1/pic_101.jpg
/usr/res/boot-play/part0/pic_023.jpg
/sbin/music-play
/usr/res/boot-play/part1/pic_144.jpg
/usr/res/boot-play/part1/pic_109.jpg
/usr/res/boot-play/part0/pic_066.jpg
/usr/res/boot-play/part1/pic_187.jpg
/usr/res/boot-play/part1/pic_110.jpg
/usr/res/boot-play/part0/pic_032.jpg
/usr/res/boot-play/part1/pic_153.jpg
/usr/res/boot-play/part1/pic_118.jpg
/usr/res/boot-play/part0/pic_075.jpg
/usr/res/boot-play/part1/pic_196.jpg
/usr/res/boot-play/part1/pic_200.jpg
/usr/res/boot-play/part0/pic_041.jpg
/usr/res/boot-play/part0/pic_006.jpg
/usr/res/boot-play/part1/pic_162.jpg
/usr/res/boot-play/part1/pic_127.jpg
/usr/res/boot-play/part1/pic_208.jpg
/usr/res/boot-play/part0/pic_084.jpg
/usr/res/boot-play/part0/pic_049.jpg
/usr/res/boot-play/part1/pic_089.jpg
/usr/res/boot-play/part0/pic_050.jpg
/usr/res/boot-play/part0/pic_015.jpg
/usr/res/boot-play/part1/pic_090.jpg
/usr/res/boot-play/part1/pic_171.jpg
/usr/res/boot-play/part1/pic_136.jpg
/usr/res/boot-play/part0/pic_058.jpg
/usr/res/boot-play/part1/pic_098.jpg
/usr/res/boot-play/part1/pic_179.jpg
/usr/res/boot-play/part1/pic_102.jpg
/usr/res/boot-play/part0/pic_024.jpg
/usr/res/boot-play/part1/pic_180.jpg
/usr/res/boot-play/part1/pic_145.jpg
/usr/res/boot-play/part0/pic_067.jpg
/usr/res/boot-play/part1/pic_188.jpg
/usr/res/boot-play/part1/pic_111.jpg
/usr/res/boot-play/part0/pic_033.jpg
/usr/res/boot-play/part1/pic_154.jpg
/usr/res/boot-play/part1/pic_119.jpg
/usr/res/boot-play/part0/pic_076.jpg
/usr/res/boot-play/part1/pic_197.jpg
/etc/init.d/play
/usr/res/boot-play/part1/pic_120.jpg
/usr/res/boot-play/part1/pic_201.jpg
/usr/res/boot-play/part0/pic_042.jpg
/usr/res/boot-play/part0/pic_007.jpg
/usr/res/boot-play/part1/pic_163.jpg
/usr/res/boot-play/part1/pic_128.jpg
/usr/res/boot-play/part1/pic_209.jpg
/usr/res/boot-play/part0/pic_085.jpg
/usr/res/boot-play/part1/pic_210.jpg
/usr/res/boot-play/part0/pic_051.jpg
/usr/res/boot-play/part0/pic_016.jpg
/usr/res/boot-play/part1/pic_091.jpg
/usr/res/boot-play/part1/pic_172.jpg
/usr/res/boot-play/part1/pic_137.jpg
/usr/res/boot-play/part0/pic_059.jpg
/usr/res/boot-play/part1/pic_099.jpg
/usr/res/boot-play/part1/pic_103.jpg
/usr/res/boot-play/part0/pic_060.jpg
/usr/res/boot-play/part0/pic_025.jpg
/usr/res/boot-play/part1/pic_181.jpg
/usr/res/boot-play/part1/pic_146.jpg
/usr/res/boot-play/part0/pic_068.jpg
/usr/res/boot-play/part1/pic_189.jpg
/usr/res/boot-play/part1/pic_112.jpg
/usr/res/boot-play/part0/pic_034.jpg
/usr/res/boot-play/part1/pic_190.jpg
/usr/res/boot-play/part1/pic_155.jpg
/usr/res/boot-play/part0/pic_077.jpg
/usr/res/boot-play/part1/pic_198.jpg
/usr/res/boot-play/part0/pic_000.jpg
/usr/res/boot-play/part1/pic_121.jpg
/usr/res/boot-play/part1/pic_202.jpg
/usr/res/boot-play/part0/pic_043.jpg
/usr/res/boot-play/part0/pic_008.jpg
/usr/res/boot-play/part1/pic_164.jpg
/usr/res/boot-play/part1/pic_129.jpg
/usr/res/boot-play/part0/pic_086.jpg
/usr/res/boot-play/part1/pic_130.jpg
/usr/res/boot-play/part1/pic_211.jpg
/usr/res/boot-play/part0/pic_052.jpg
/usr/res/boot-play/part0/pic_017.jpg
/usr/res/boot-play/part1/pic_092.jpg
/usr/res/boot-play/part1/pic_173.jpg
/usr/res/boot-play/part1/pic_138.jpg
/usr/res/boot-play/part1/pic_104.jpg
/usr/res/boot-play/part0/pic_061.jpg
/usr/res/boot-play/part0/pic_026.jpg
/usr/res/boot-play/part1/pic_182.jpg
/usr/res/boot-play/part1/pic_147.jpg
/usr/res/boot-play/part0/pic_069.jpg
/usr/res/boot-play/part1/pic_113.jpg
/usr/res/boot-play/part0/pic_070.jpg
/usr/res/boot-play/part0/pic_035.jpg
/usr/res/boot-play/part1/pic_191.jpg
/usr/res/boot-play/part1/pic_156.jpg
/usr/res/boot-play/part0/pic_078.jpg
/usr/res/boot-play/part1/pic_199.jpg
/usr/res/boot-play/part0/pic_001.jpg
/usr/res/boot-play/part1/pic_122.jpg
/usr/res/boot-play/part1/pic_203.jpg
/usr/res/boot-play/part0/pic_044.jpg
/usr/res/boot-play/part0/pic_009.jpg
/usr/res/boot-play/part1/pic_165.jpg
/usr/res/boot-play/part0/pic_087.jpg
/usr/res/boot-play/part0/pic_010.jpg
/usr/res/boot-play/part1/pic_131.jpg
/usr/res/boot-play/part1/pic_212.jpg
/usr/res/boot-play/part0/pic_053.jpg
/usr/res/boot-play/part0/pic_018.jpg
/usr/res/boot-play/part1/pic_093.jpg
/usr/res/boot-play/part1/pic_174.jpg
/usr/res/boot-play/part1/pic_139.jpg
/usr/res/boot-play/part1/pic_140.jpg
/usr/res/boot-play/part1/pic_105.jpg
/usr/res/boot-play/part0/pic_062.jpg
/usr/res/boot-play/part0/pic_027.jpg
/usr/res/boot-play/part1/pic_183.jpg
/usr/res/boot-play/part1/pic_148.jpg
/usr/res/boot-play/part1/pic_114.jpg
/usr/res/boot-play/part0/pic_071.jpg
/usr/res/boot-play/part0/pic_036.jpg
/usr/res/boot-play/part1/pic_192.jpg
/usr/res/boot-play/part1/pic_157.jpg
/usr/res/boot-play/part0/pic_079.jpg
/usr/res/boot-play/part0/pic_002.jpg
/usr/res/boot-play/part1/pic_123.jpg
/usr/res/boot-play/part1/pic_204.jpg
/usr/res/boot-play/part0/pic_080.jpg
/usr/res/boot-play/part0/pic_045.jpg
/usr/res/boot-play/part1/pic_166.jpg
/usr/res/boot-play/part0/pic_088.jpg
